# pkg/raysl-desktop/raysl-desktop.install
post_install() {
	systemctl --user enable --now redshift.service
  post_upgrade
}

post_upgrade() {

}
